# transcription-offline.service
[Unit]
Description=Desktop Transcription Tool (Offline)
Documentation=https://github.com/rokytnice/desktop_transcription_tool
After=network.target sound.target multi-user.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=aroc
Group=aroc
Environment="AUDIO_DEVICE=16"
Environment="AUDIO_OUTPUT_DEVICE=16"
Environment="WHISPER_MODEL=small"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/aroc/.Xauthority"
WorkingDirectory=/home/aroc/projects/desktop_transcription_tool/offline
ExecStart=/home/aroc/projects/desktop_transcription_tool/offline/.venv/bin/python /home/aroc/projects/desktop_transcription_tool/offline/transcription_offline.py
Restart=on-abort
RestartSec=30
StandardOutput=journal
StandardError=journal

# Allow access to input devices and audio
ProtectSystem=strict
ProtectHome=no
NoNewPrivileges=false

[Install]
WantedBy=graphical.target
