# transcription-offline.service
[Unit]
Description=Desktop Transcription Tool (Offline)
Documentation=https://github.com/rokytnice/desktop_transcription_tool
After=pipewire.service pipewire-pulse.service

[Service]
Type=simple
Environment="AUDIO_DEVICE=7"
Environment="AUDIO_OUTPUT_DEVICE=7"
Environment="WHISPER_MODEL=small"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
WorkingDirectory=/home/aroc/projects/desktop_transcription_tool/offline
ExecStart=/home/aroc/projects/desktop_transcription_tool/offline/.venv/bin/python /home/aroc/projects/desktop_transcription_tool/offline/transcription_offline.py
Restart=on-abort
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
